# service/phishdefender.service
[Unit]
Description=PhishDefender - Phishing Link Detection Service
After=network.target

[Service]
ExecStart=/usr/bin/env python3 /usr/local/bin/phishdefender
Restart=always
User=phishdefender
Group=phishdefender

[Install]
WantedBy=multi-user.target
